# scripts/packaging/cat-detector.service
[Unit]
Description=Cat Detector - Webcam-based cat detection with notifications
After=network.target

[Service]
Type=simple
User=griswald
Environment=ORT_DYLIB_PATH=/opt/cat-detector/lib/libonnxruntime.so
ExecStart=/usr/local/bin/cat-detector run --config /etc/cat-detector/config.toml
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
